# Centos6/smcactivemq.conf
### /etc/init/

respawn
start on runlevel [345]
exec runuser - smc -c "/opt/smc/mq/activemq/bin/activemq console"

