# Centos7/redis.service.d/override.conf
[Service]
Type=simple
ExecStart=
ExecStart=/usr/bin/redis-server /etc/redis.conf