# Centos7/elasticsearch.service.d/override.conf
[Unit]
Before=smceventmonitor.service
Before=smcfilestore.service
[Service]
TimeoutStartSec=240
LimitMEMLOCK=infinity
