# Centos7/getsyserr.timer
[Unit]
Description=Get System Error Timer

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
